# initrd_source/skel/etc/udev/rules.d/80-programs.rules
# firmware class requests
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/lib/udev/firmware.sh"

