# debian/systemimager-server.conffiles
/etc/systemimager/systemimager.conf
/etc/systemimager/rsyncd.conf
/etc/init.d/systemimager
