# debian/systemimager-server.conffiles
/etc/systemimager/server.conf
/etc/systemimager/rsyncd.conf
/etc/init.d/systemimager
