# debian/va-systemimager.postrm
if [ purge = "$1" ]; then
update-rc.d rsync remove >/dev/null
fi
